# middleman.service.in
[Unit]
Description=Helium radio adjustment middleman

[Service]
ExecStart=@@INSTALLDIR@@/run.sh
WorkingDirectory=@@INSTALLDIR@@

[Install]
WantedBy=multi-user.target
