# middleman.service.in
[Unit]
Description=Helium radio adjustment middleman

[Service]
ExecStart=@@DESTROOT@@/run.sh
WorkingDirectory=@@DESTROOT@@
Restart=always

[Install]
WantedBy=multi-user.target
